# docker/root/etc/services.d/themepark/run
#!/usr/bin/with-contenv bash

  exec \
        python3 /config/www/themes.py