# docker-mods/qbittorrent/root/etc/s6-overlay/s6-rc.d/init-mod-themepark-add-package/run
#!/usr/bin/with-contenv bash 

if ! [[ -x "$(command -v git)" ]]; then
    echo "git" >> /mod-repo-packages-to-install.list
    echo "perl" >> /mod-repo-packages-to-install.list
fi
