# docker-mods/vuetorrent/root/etc/s6-overlay/s6-rc.d/init-mod-themepark-add-package/run
#!/usr/bin/with-contenv bash 

if ! [[ -x "$(command -v jq)" ]]; then
    echo "jq" >> /mod-repo-packages-to-install.list
fi
